# config/systemd/user/speak-to-ai.service
[Unit]
Description=Speak to AI - Offline Speech-to-Text Daemon
Documentation=https://github.com/speak-to-ai/speak-to-ai
After=pipewire.service pipewire-pulse.service
Requires=pipewire.service

[Service]
Type=simple
ExecStart=/usr/bin/speak-to-ai -config %h/.config/speak-to-ai/config.yaml
Restart=on-failure
RestartSec=5

# Environment
Environment="XDG_RUNTIME_DIR=/run/user/%U"
Environment="PULSE_SERVER=unix:/run/user/%U/pulse/native"
Environment="PIPEWIRE_RUNTIME_DIR=/run/user/%U"

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
